# kitchensync.service
[Unit]
Description=KitchenSync Auto-Start Service
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=kitchensync
Group=kitchensync
WorkingDirectory=/home/kitchensync/kitchenSync
ExecStart=/usr/bin/python3 /home/kitchensync/kitchenSync/kitchensync.py
Restart=always
RestartSec=10
Environment=DISPLAY=:0

# Wait for desktop to be ready
ExecStartPre=/bin/sleep 30

[Install]
WantedBy=graphical.target
