# kitchensync.service
[Unit]
Description=KitchenSync Auto-Start Service
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=kitchensync
Group=kitchensync
WorkingDirectory=/home/kitchensync/kitchenSync
ExecStart=/home/kitchensync/kitchenSync/start_clean.sh
Restart=always
RestartSec=10

# Wait for desktop to be ready
ExecStartPre=/bin/sleep 30

[Install]
WantedBy=default.target
